# deploy/systemd/oil-research.service
[Unit]
Description=Zhonglu refined oil research system
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/oil-research
RemainAfterExit=yes
ExecStart=/usr/bin/docker compose -p oilresearch up -d --build
ExecStop=/usr/bin/docker compose -p oilresearch down
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
